# docker/20210707/prod/docker-compose/01-webapp/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true

[program:php-fpm]
command=/usr/sbin/php-fpm -c /etc/php-fpm.conf
autostart=true